# windows/task-scheduler.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Gaming Cafe - Task Scheduler Definition
  
  To import this manually:
  1. Open Task Scheduler (search in Start menu)
  2. Click "Import Task" in the right panel
  3. Select this file
  4. Change <UserId> to your Windows username
  5. Click OK
  
  OR just run setup.bat as Administrator - it does this automatically.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Starts Gaming Cafe backend and frontend servers on login</Description>
    <URI>\GamingCafe</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <!-- Change this to your Windows username -->
      <UserId>YOUR_WINDOWS_USERNAME</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\gc-be\windows\start-cafe.bat</Command>
    </Exec>
  </Actions>
</Task>
